# scripts/systemd/hyxlab-backup.service
[Unit]
Description=hyxlab rotated DuckDB backups (set HYXLAB_BACKUP_DIR for off-box)
[Service]
Type=oneshot
OOMScoreAdjust=500
WorkingDirectory=/home/devs/workspace/hyxrestration-stable
ExecStart=/home/devs/workspace/hyxrestration-stable/.venv/bin/python -u -m collector.backup
