# scripts/systemd/hyxlab-signals.service
[Unit]
Description=hyxlab daily signal pull (ALFRED econ vintages + GDELT news)
[Service]
Type=oneshot
WorkingDirectory=/home/devs/workspace/hyxrestration-stable
ExecStart=/home/devs/workspace/hyxrestration-stable/.venv/bin/python -u -m collector.signals
