# scripts/systemd/hyxlab-simui.service
[Unit]
Description=hyxlab simui market-replay terminal (localhost:8877, paper only)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/devs/workspace/hyxrestration-stable
ExecStart=/home/devs/workspace/hyxrestration-stable/.venv/bin/python -u -m simulator.simui
Restart=always
RestartSec=10
MemoryMax=1G

[Install]
WantedBy=default.target
